# cookbooks/dmon/templates/default/start_frontend.service.erb
[Unit]
Description=<%=@componentName%> service
[Service]
ExecStart=<%=@home%>/bin/start.sh
ExecStop=<%=@home%>/bin/stop.sh
User=root
Type=simple
[Install]
WantedBy=multi-user.target
